# socket-proxy/supervisor/conf.d/backend.conf
[program:backend]
command=npm run serve
directory=/app
