# socket-proxy/supervisor/conf.d/backend.conf
[program:backend]
command=npm run dev
# directory=/tmp
